# provisioning/systemd/backup-horse-racing.timer
[Unit]
Description=Run backup-horse-racing.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
