# provisioning/systemd/crawl-github-activity.timer
[Unit]
Description=Run crawl-github-activity.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
