# provisioning/systemd/crawl-redmine-activity.timer
[Unit]
Description=Run crawl-redmine-activity.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
