# provisioning/systemd/update-mydns.service
[Unit]
Description=Update mydns ip address

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-mydns.sh

[Install]
WantedBy=multi-user.target
